# e/systemd/app@.service
[Service]
SyslogIdentifier=%i
ExecStart=/usr/bin/env %i
