# e/systemd/conky.service
[Unit]
After=aura-oneshot.service

[Service]
ExecStartPre=/bin/sleep 1
ExecStart=/usr/bin/env conky
