# e/systemd/dbus.socket
[Socket]
ListenStream=/run/user/1000/session_bus_socket
