# e/systemd/enlightenment.service
[Unit]
Before=aura-oneshot.service xsettings.service

[Service]
SyslogIdentifier=%p
KillMode=none
ExecStart=/usr/bin/env enlightenment_start
