# e/systemd/enlightenment.service
[Unit]
Before=aura-oneshot.service xinput.service

[Service]
ExecStart=/usr/bin/enlightenment_start
ExecStopPost=-/usr/bin/systemctl --fail --no-block --user exit
