# e/systemd/yeahconsole.service
[Service]
ExecStart=/usr/bin/env yeahconsole
