# e/systemd/yeahconsole.service
[Service]
SyslogIdentifier=%p
ExecStart=/usr/bin/env yeahconsole
