# systemd/user/aqm-alerts.socket
[Socket]
ListenDatagram=0.0.0.0:8153

[Install]
WantedBy=sockets.target
