# systemd/user/aqm-alerts.service
[Service]
Type=exec
ExecStart=env notify.aqm -z 127.0.0.1:8158 \
	-U critical -t0 -i aqm --udp-snooze 2h --exit-inactivity-timeout 15m
