# systemd/user/autocutsel.service
[Unit]
After=xorg.target
PartOf=startx.target

[Service]
SyslogIdentifier=%p
Restart=on-failure
ExecStart=/usr/bin/env autocutsel
