# systemd/user/conky.service
[Unit]
After=xorg.target background.target
PartOf=startx.target

[Service]
SyslogIdentifier=%p
Restart=on-abort
ExecStart=env conky -q
