# systemd/user/cron-check-dev-nodes.service
[Unit]
PartOf=notification-thing.service

[Service]
Type=oneshot
ExecStart=env cron.check.dev-nodes
