# systemd/user/feh-screen.service
[Service]
Type=oneshot
StandardInput=socket
StandardOutput=journal
ExecStart=feh-screen -F=-g=${VO_AUX_POS}
