# systemd/user/feh-screen.socket
[Socket]
SocketMode=0600
ListenFIFO=%h/.config/feh/screen.sock

[Install]
WantedBy=sockets.target
