# systemd/user/mpv-urls.socket
[Socket]
SocketMode=0600
ListenFIFO=%h/.mpv/urls.sock

[Install]
WantedBy=sockets.target
