# systemd/user/nodpi-https-proxy.service
[Service]
Type=exec
ExecStart=nodpi-https-proxy -t 600
