# systemd/user/notify-power.service
[Service]
ExecStart=/usr/local/bin/notify.power
