# systemd/user/riet.timer
[Unit]
After=xorg.target background.target notification-thing.service
PartOf=startx.target

[Timer]
OnActiveSec=30
