# systemd/user/yeahconsole.timer
[Unit]
After=xorg.target background.target
PartOf=startx.target

[Timer]
OnActiveSec=18
