# e/systemd/user/pulseaudio.service.d/manual.conf
[Service]
ExecStart=
ExecStart=/usr/bin/pulseaudio --daemonize=no -nF /etc/pulse/manual.pa
