# e/systemd/user/app@.service
[Unit]
After=xorg.target
PartOf=startx.target

[Service]
SyslogIdentifier=%i
ExecStart=/usr/bin/env %i
