# e/systemd/user/app@.service
[Unit]
After=X.target

[Service]
SyslogIdentifier=%i
ExecStart=/usr/bin/env %i
