# e/systemd/user/background.target
[Unit]
Description=X background image, kinda DE-init point
Requires=X.target
After=X.target
