# e/systemd/user/background.target
[Unit]
Description=X background image, kinda DE-init point
After=X.target
