# e/systemd/user/conky.service
[Unit]
After=X.target aura-oneshot.service

[Service]
SyslogIdentifier=%p
ExecStart=/usr/bin/env conky -q
