# e/systemd/user/dbus.socket
[Socket]
ListenStream=/run/user/fraggod/dbus/user_bus_socket
