# e/systemd/user/X-start-delay.service
[Unit]
After=X.service
Before=X.target

[Service]
Type=oneshot
ExecStart=/bin/sleep 2
