# e/systemd/user/X.service
[Unit]
Requires=X.target
Before=X.target

[Service]
SyslogIdentifier=%p
ExecStart=/usr/lib/systemd/systemd-multi-seat-x $DISPLAY
