# e/systemd/user/xorg.target
[Unit]
Description=X server
PartOf=startx.target
