# e/systemd/user/xorg-start-delay.service
[Unit]
DefaultDependencies=no
After=xorg@1.service 

[Service]
ExecStart=/usr/bin/sleep 2
