# e/systemd/user/xsettings.timer
[Unit]
After=xorg.target
PartOf=startx.target

[Timer]
OnActiveSec=20
