# e/systemd/user/yeahconsole.service
[Unit]
After=X.target background.target
PartOf=startx.target

[Service]
SyslogIdentifier=%p
ExecStart=/usr/bin/env yeahconsole
